# e-paper_monitor.service
[Unit]
Description=e-Paper System Monitor (WaveShare 2.13" HAT+)
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/e-paper_monitor.py
WorkingDirectory=/usr/local/bin
Restart=on-failure
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=e-paper-monitor

[Install]
WantedBy=multi-user.target
